# packages/uniflow/root/Applications/UniFlow/CanonPS.mpkg/Contents/Packages/Canon PS Installer.pkg/Contents/Resources/Software/OSDiffs/launchd/jp.co.canon.CUPSPS2.BG.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>jp.co.canon.CUPSPS2.BackGrounder</string>
	<key>KeepAlive</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Printers/Canon/CUPSPS2/BackGrounder/Canon PS2 BackGrounder.app/Contents/MacOS/Canon PS2 BackGrounder</string>
	</array>
</dict>
</plist>
